# infra/systemd/bi-rmp-staging.service.example
[Unit]
Description=BI-RMP Staging Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/harcker8119/BI-RMP-STAGING
EnvironmentFile=/home/harcker8119/BI-RMP-STAGING/.env.staging.runtime
ExecStart=/home/harcker8119/BI-RMP-STAGING/.venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8101 --app-dir Backend
Restart=on-failure
RestartSec=5
User=harcker8119
Group=harcker8119
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
