# infra/systemd/bi-rmp-staging.service.example
[Unit]
Description=BI-RMP Staging Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__STAGING_APP_DIR__
EnvironmentFile=__STAGING_ENV_FILE__
ExecStart=__STAGING_APP_DIR__/.venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8101 --app-dir Backend
Restart=on-failure
RestartSec=5
User=__STAGING_USER__
Group=__STAGING_GROUP__
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
